# overlay/etc/rc.d/mangadexathome
#!/bin/sh

# PROVIDE: mangadexathome

. /etc/rc.subr

name="mangadexathome"
user="mangadexathome"
rcvar=`set_rcvar`
pidfile="/var/run/mangadexathome"
start_cmd="${name}_start"

load_rc_config $name

mangadexathome_start()
{
    if checkyesno ${rcvar}; then
      cd /usr/home/mangadexathome/mangadexathome
      /usr/sbin/daemon p $pidfile -u $user java -Dfile-level=info -Dstdout-level=off -jar mangadex_at_home.jar > /dev/null 2>&1
    fi
}

run_rc_command "$1"

